# deploy/alibaba-ecs/systemd/yiting-recorder-heartbeat.service
[Unit]
Description=YITING Recorder Heartbeat
After=yiting-gateway.service
Requires=yiting-gateway.service

[Service]
Type=simple
WorkingDirectory=/opt/yiting
EnvironmentFile=/opt/yiting/.env
ExecStart=/opt/yiting/.venv/bin/python -m agents.recorder.heartbeat
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
