# deploy/alibaba-ecs/systemd/yiting-victim.service
[Unit]
Description=YITING Victim App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/yiting
EnvironmentFile=/opt/yiting/.env
ExecStart=/opt/yiting/.venv/bin/uvicorn app:app --app-dir victim-app --host 127.0.0.1 --port 9000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
